# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:proxy-server]
command=node --max-old-space-size=2048 --expose-gc /app/server/proxy-server.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/proxy-server.err.log
stdout_logfile=/var/log/supervisor/proxy-server.out.log
environment=NODE_ENV="production",PORT="3001",PDF_SERVER_INTERNAL_TOKEN="%(ENV_PDF_SERVER_INTERNAL_TOKEN)s",JWT_SECRET="%(ENV_JWT_SECRET)s",CSRF_SECRET="%(ENV_CSRF_SECRET)s",CACHE_BACKEND="%(ENV_CACHE_BACKEND)s",CACHE_REDIS_URL="%(ENV_CACHE_REDIS_URL)s",CACHE_KEY_PREFIX="%(ENV_CACHE_KEY_PREFIX)s"
priority=5

[program:pdf-server]
command=node /app/pdf-server/server.cjs
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/pdf-server.err.log
stdout_logfile=/var/log/supervisor/pdf-server.out.log
environment=NODE_ENV="production",PDF_SERVER_PORT="3002",PUPPETEER_EXECUTABLE_PATH="/usr/bin/google-chrome-stable",PDF_SERVER_INTERNAL_TOKEN="%(ENV_PDF_SERVER_INTERNAL_TOKEN)s",JWT_SECRET="%(ENV_JWT_SECRET)s",CSRF_SECRET="%(ENV_CSRF_SECRET)s"
priority=10


